# services.d/puma/run
#!/command/with-contenv sh

cd /app
rails s
echo "puma running"